# Dotfiles/Wallpapers/dynamic-wallpaper.timer
[Unit]
Description=Run the dynamic wallpaper script each hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
