# Dotfiles/General/system-backup.service
[Unit]
Description=Perform a system backup

[Service]
Type=oneshot
ExecStart=system-backup --create --scheduled

[Install]
WantedBy=default.target
